# roles/jhaycr-local.safe_shutdown/templates/safe-shutdown.service.j2
[Unit]
Description=Safe System Shutdown Handler
DefaultDependencies=no
Before=shutdown.target reboot.target halt.target
Before=umount.target
After=docker.service

[Service]
Type=oneshot
Environment=SYSTEMD_SHUTDOWN=yes
ExecStart=/usr/local/sbin/safe-shutdown.sh
TimeoutStartSec=600
RemainAfterExit=yes

[Install]
WantedBy=shutdown.target reboot.target halt.target